# supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null

[program:nginx]
user=root
command=nginx -g 'daemon off;'
autorestart=true

[program:crond]
user=root
command=/usr/sbin/cron
autostart=true
